# Software/scripts/lowpowerscripts_raspberryconnect/lowpower.service
[Unit]
Description="Powers Down some services after 10 minutes"

[Service]
ExecStart=/usr/bin/lowpower.sh
