# Firmware/5.x/scripts/MothPower/lowpower.timer
[Unit]
Description="for low power mode after 10 minutes"

[Timer]
OnBootSec=10min

[Install]
WantedBy=multi-user.target
